# lidarr/root/etc/services.d/lidarr/run
#!/usr/bin/env bashio

# cd /app/lidarr/bin || exit



exec \
    s6-setuidgid abc /app/lidarr/bin/Lidarr \
    -nobrowser -data=$(bashio::config 'Config_Folder')
